# config/launchd/com.stock.fund-tail-advice.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.stock.fund-tail-advice</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Volumes/hw-val/Dev/stock/scripts/run_daily_fund_tail_advice.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Volumes/hw-val/Dev/stock</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>50</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>50</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>50</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>50</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>50</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Volumes/hw-val/Dev/stock/reports/fund_tail_advice/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Volumes/hw-val/Dev/stock/reports/fund_tail_advice/launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
